# files/etc/uci-defaults/99-qwrt-defaults
#!/bin/sh
# 99-qwrt-defaults — 首启动一次性配置
#
# 当前策略：尽量保持默认系统行为，不主动接管 DNS，
# 只处理服务默认状态，避免用户首次开机就被复杂代理栈卡住。

set -e

# sing-box：默认 disabled（用户填入真实节点后再启用）
[ -x /etc/init.d/sing-box ] && /etc/init.d/sing-box disable

# easytier：默认 disabled（用户填入网络密钥后再启用）
[ -x /etc/init.d/easytier ] && /etc/init.d/easytier disable

# 若历史残留服务存在，也统一关掉，保持 sing-box 裸核运行
[ -x /etc/init.d/momo ] && /etc/init.d/momo disable
[ -x /etc/init.d/mosdns ] && /etc/init.d/mosdns disable
[ -x /etc/init.d/daed ] && /etc/init.d/daed disable

uci commit
exit 0
